# service/brtg-pk-service/log-backup.service
[Unit]
Description=Log Rotate Service

[Service]
ExecStart=/usr/local/bin/log_backup.sh
